# commotionfeed/commotionbase/files/etc/hotplug.d/iface/00-services_predown
[ "$ACTION" = "ifdown" ] && [ "$PROTO" = "meshif" -o "$PROTO" = "apif" -o "$PROTO" = "plugif ] && {
  ACTION="predown" /sbin/hotplug-call "services" &
}
