# commotionfeed/commotionbase/files/etc/hotplug.d/iface/99-services_postup
[ "$ACTION" = "ifup" ] && [ "$PROTO" = "meshif" -o "$PROTO" = "apif" -o "$PROTO" = "plugif" ] && {
  ACTION="postup" /sbin/hotplug-call "services" &
}
