# commotionfeed/commotionbase/files/etc/hotplug.d/services/99-firewall
case $PROTO in
  meshif)
      case $ACTION in
        postup)
          /etc/init.d/firewall restart
        ;;
      esac
  ;;
esac
