# default-files/etc/hotplug.d/iface/90-thisnode
#!/bin/sh
#DEBUG="echo"

. /lib/functions.sh
. /lib/functions/commotion.sh
. /lib/config/uci.sh
. /lib/firewall/core_interface.sh

if [ "$ACTION" == "ifup" ] && [ $(/usr/bin/commotion state "$DEVICE" mode) == "wired" ]; then
	localip=$(/usr/bin/commotion localipset)
        THISNODE=$localip" thisnode"                    
        if grep -Fq thisnode /etc/hosts; then                                                
                sed -ie "s/^.*thisnode$/$THISNODE/" /etc/hosts || \      
                        logger -t commotion.hotplug.thisnode -s "Couldn't edit /etc/hosts"
                logger -t commotion.hotplug.thisnode -s "Updating IP address for 'thisnode' host alias"
        else                                                                                           
                echo $THISNODE >> /etc/hosts                                                           
                logger -t commotion.hotplug.thisnode -s "Creating entry for 'thisnode' host alias"

	fi 
	
	ALIAS="thisnode"
	
	logger -t commotion.hotplug.thisnode -s "Checking UCI entry for $ALIAS alias"
	
	if [ grep -Fq $ALIAS /etc/config/network] && [grep -Fq $localip ]; then
	      logger -t commotion.hotplug.thisnode -s "$ALIAS alias exists with ip address $localip"
	
	else 
		uci set network.$ALIAS=alias
		logger -t commotion.hotplug.thisnode -s "Creating UCI entry for '$ALIAS' alias"
		uci set network.$ALIAS.interface=plug
		uci set network.$ALIAS.proto=static
		
		uci set network.$ALIAS.ipaddr=$localip
		logger -t commotion.hotplug.thisnode -s "Setting '$ALIAS' ip address to $localip"
		uci set network.$ALIAS.netmask=255.255.0.0
		
		uci commit network
	fi	
		#NOTE: this method of configuring aliases is deprecated in more recent versions of openwrt
		
                                                                                     
	kill -s HUP `pgrep dnsmasq`
fi