# default-files/etc/hotplug.d/iface/90-thisnode
#!/bin/sh
#DEBUG="echo"

. /lib/functions.sh
. /lib/functions/commotion.sh
. /lib/config/uci.sh
. /lib/firewall/core_interface.sh

if [ "$ACTION" == "ifup" ] && [ $(/usr/bin/commotion state "$DEVICE" mode) == "wired" ]; then
	localip=$(/usr/bin/commotion localipset)
        THISNODE=$localip" thisnode"                    
        if grep -Fq thisnode /etc/hosts; then                                                
                sed -ie "s/^.*thisnode$/$THISNODE/" /etc/hosts || \      
                        logger -t commotion.hotplug.thisnode -s "Couldn't edit /etc/hosts"
                logger -t commotion.hotplug.thisnode -s "Updating IP address for 'thisnode' host alias"
        else                                                                                           
                echo $THISNODE >> /etc/hosts                                                           
                logger -t commotion.hotplug.thisnode -s "Creating entry for 'thisnode' host alias"

	fi 
		ALIAS="thisnode"
				
		uci set network.$ALIAS=alias
		logger -t commotion.hotplug.thisnode -s "Creating uci entry for 'thisnode' alias"
		uci set network.$ALIAS.interface=plug
		uci set network.$ALIAS.proto=static
		
		uci set network.$ALIAS.ipaddr=$localip
		logger -t commotion.hotplug.thisnode -s "Setting 'thisnode' ip address to $localip"
		uci set network.$ALIAS.netmask=255.255.0.0
		
                                                                                     
	kill -s HUP `pgrep dnsmasq`
fi